# profile.d/mongo.sh
export MONGO_HOME=/data/lnmp/mongo
export PATH=$MONGO_HOME/bin:$PATH
