# profile.d/mysql.sh
export MYSQL_HOME=/data/lnmp/mysql
export PATH=$MYSQL_HOME/bin:$PATH
