# profile.d/php.sh
export PHP_HOME=/data/lnmp/php56
export PATH=$PHP_HOME/bin:$PATH
