# profile.d/redis.sh
export REDIS_HOME=/data/lnmp/redis
export PATH=$REDIS_HOME/bin:$PATH
